# etc/cron.daily/cyr_expire
#!/usr/bin/sh

cyr_expire -X 7 -v

exit 0
