# etc/systemd/system/clamd@amavisd.service.d/override.conf
[Service]
ExecStartPost=chmod 770 /var/run/clamd.scan
